# boot_daemon.service
[Unit]
Description=BOOT net modules
Wants=multi-user.target

[Service]
ExecStart=/usr/local/tools/boot.sh

[Install]
WantedBy=multi-user.target